# remote-mouse.service
[Unit]
Description=Remote Mouse Driver
After=multi-user.target

[Service]
Type=simple
User=root
ExecStart=/path/to/your/venv/bin/remote-mouse-driver
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target